# systemd/gloaming-brightness-schedule.service
[Unit]
Description=Gloaming: fade monitor brightness to the current schedule target

[Service]
Type=oneshot
ExecStart=%h/Projects/gloaming/brightness/schedule-brightness.sh
